# gitland-bot.service
[Unit]
Description=Run wmwws Gitland bot

[Service]
ExecStart=/usr/bin/gitland-bot

[Install]
WantedBy=multi-user.target
